# etc/supervisord.d/httpd.ini
[program:httpd]
command=httpd -X
startsecs=2
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
