# etc/supervisord.d/icinga2.ini
[program:icinga2]
command=/usr/sbin/icinga2 --config /etc/icinga2/icinga2.conf
priority=1
startsecs=5
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
