# etc/supervisord.d/mysqld.ini
[program:mysqld]
command=/root/bin/start_mysqld.sh
priority=1
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

