# etc/supervisord.d/setup_icinga2web.ini
[program:setup_icinga2web]
command=/root/bin/setup_icinga2web.sh
startsec=0
autoretries=0
priority=1
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

